# heinz.service
[Unit]
Description=Heinzelmaennekes home automation

[Service]
ExecStart=/usr/bin/heinzd
Type=simple

[Install]
WantedBy=multi-user.target
